# deployment/gunicorn.service
[Unit]
Description=Gunicorn daemon for ShadCoding Task1 Django Application
After=network.target

[Service]
Type=notify
# Run as deploy user
User=deploy
Group=deploy

# Working directory
WorkingDirectory=/home/deploy/shadcoding-task1/backend

# Environment variables
Environment="PATH=/home/deploy/shadcoding-task1/venv/bin"
EnvironmentFile=/home/deploy/shadcoding-task1/backend/.env

# Gunicorn command
ExecStart=/home/deploy/shadcoding-task1/venv/bin/gunicorn \
    --config /home/deploy/shadcoding-task1/deployment/gunicorn.conf.py \
    backend.wsgi:application

# Restart policy
Restart=on-failure
RestartSec=5s

# Resource limits
LimitNOFILE=4096

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
